# kit.cron.d
## Adjust this line...
MAILTO=mkronvold

0 * * * * mkronvold /e2open/home/mkronvold/src/kit/push.sh > /dev/null
